# scripts/com.jibun_os.notion-coordinate-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  コーディネート予定をカレンダー（西岡慎也）→ Notion タスク管理 DB へ毎朝 8:01 に同期

  初回のみ（パスは環境に合わせて修正）:
    cp scripts/com.jibun_os.notion-coordinate-sync.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.jibun_os.notion-coordinate-sync.plist

  止める:
    launchctl unload ~/Library/LaunchAgents/com.jibun_os.notion-coordinate-sync.plist

  ログ:
    tasks/.gcal-oauth/sync-notion.log
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.jibun_os.notion-coordinate-sync</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/hirotaseiji/Desktop/自分OS/scripts/run-coordinate-sync-notion.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>8</integer>
    <key>Minute</key>
    <integer>1</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/hirotaseiji/Desktop/自分OS/tasks/.gcal-oauth/launchd-stdout-notion.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/hirotaseiji/Desktop/自分OS/tasks/.gcal-oauth/launchd-stderr-notion.log</string>
</dict>
</plist>
